# dashcam.service
[Unit]
Description=Raspberry Pi Dashcam System
After=network.target
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dashcam
ExecStart=/home/pi/dashcam/.venv/bin/python /home/pi/dashcam/dashcam.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 